# .worktrees/pr325/systemd/cryptotrader-dossier.service
[Unit]
Description=CryptoTrader Daily Dossier Generator
After=network.target postgresql.service cryptotrader-api.service

[Service]
Type=oneshot
WorkingDirectory=%h/cryptotrader
ExecStart=%h/cryptotrader/.venv/bin/python -m scripts.generate_dossiers --exchange bitfinex --delay 15
Environment=PYTHONPATH=%h/cryptotrader
EnvironmentFile=%h/cryptotrader/.env

[Install]
WantedBy=default.target
