# .worktrees/pr325/systemd/cryptotrader-dossier.timer
[Unit]
Description=CryptoTrader Daily Dossier Timer (runs at 08:00 UTC)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
