# systemd/cryptotrader-approve-workflows.service
[Unit]
Description=Auto-approve GitHub Actions workflow runs when Copilot finishes
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/cryptotrader/.venv/bin/python -m scripts.approve_workflows --daemon --interval 120
WorkingDirectory=%h/cryptotrader
Restart=on-failure
RestartSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=default.target
