# systemd/cryptotrader-frontend.service
[Unit]
Description=cryptotrader frontend (Vite dev)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/cryptotrader/frontend
Environment=NODE_ENV=development

# Run Vite dev server with proxy support
ExecStart=/usr/bin/env npm run dev -- --host 0.0.0.0 --port 5176 --strictPort

Restart=on-failure
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/cryptotrader/frontend

[Install]
WantedBy=default.target
