# systemd/cryptotrader-dossier.service
[Unit]
Description=CryptoTrader Daily Dossier Generator
After=network.target postgresql.service cryptotrader-api.service

[Service]
Type=oneshot
WorkingDirectory=/home/flip/cryptotrader
ExecStart=/home/flip/cryptotrader/.venv/bin/python -m scripts.generate_dossiers --exchange bitfinex
Environment=PYTHONPATH=/home/flip/cryptotrader
EnvironmentFile=/home/flip/cryptotrader/.env

[Install]
WantedBy=default.target
