# deployment/systemd/copilot/ct-frontend-copilot.service
[Unit]
Description=CryptoTrader frontend preview (copilot stack)
After=network-online.target ct-backend-copilot.service
Wants=network-online.target

[Service]
Type=simple
User=flip
Group=flip
WorkingDirectory=/home/flip/cryptotrader_copilot
EnvironmentFile=/home/flip/cryptotrader_copilot/.env
Environment=PORT=50000
Environment=FRONTEND_PORT=50176
Environment=VITE_API_PROXY_TARGET=http://127.0.0.1:50000
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/env npm --prefix /home/flip/cryptotrader_copilot/frontend run build
ExecStart=/usr/bin/env bash -lc 'exec npm --prefix /home/flip/cryptotrader_copilot/frontend run preview -- --host 0.0.0.0 --port "${FRONTEND_PORT}" --strictPort'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/flip/cryptotrader_copilot/frontend

[Install]
WantedBy=multi-user.target