# deployment/systemd/hermes/ct-backend-hermes.service
[Unit]
Description=CryptoTrader backend API (hermes stack)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flip
Group=flip
WorkingDirectory=/home/flip/cryptotrader_hermes
EnvironmentFile=/home/flip/cryptotrader_hermes/.env
Environment=PYTHONPATH=/home/flip/cryptotrader_hermes
Environment=PORT=51000
ExecStart=/usr/bin/env bash -lc 'exec /home/flip/cryptotrader_hermes/.venv/bin/python /home/flip/cryptotrader_hermes/scripts/run_api.py --host 0.0.0.0 --port "${PORT}"'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/flip/cryptotrader_hermes

[Install]
WantedBy=multi-user.target