# deployment/systemd/hermes/ct-legacy-hermes.service
[Unit]
Description=CryptoTrader legacy dashboard helper (hermes stack)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flip
Group=flip
WorkingDirectory=/home/flip/cryptotrader_hermes
EnvironmentFile=/home/flip/cryptotrader_hermes/.env
Environment=PYTHONPATH=/home/flip/cryptotrader_hermes
Environment=LEGACY_PORT=51787
ExecStart=/usr/bin/env bash -lc 'exec /home/flip/cryptotrader_hermes/.venv/bin/python /home/flip/cryptotrader_hermes/scripts/api_server.py --host 127.0.0.1 --port "${LEGACY_PORT}"'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/flip/cryptotrader_hermes

[Install]
WantedBy=multi-user.target